# systemd/dns-probe@.service.in
# DNS probe daemon
# Service parameter represents unique ID of DNS Probe instance.
# The ID specifies what configuration to load from YAML configuration file
# in @CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_SYSCONFDIR@/dns-probe-@BACKEND@/dns-probe.yml
[Unit]
Description=Daemon for DNS probe
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_BINDIR@/dns-probe-@BACKEND@ -n %i

[Install]
WantedBy=multi-user.target
