# systemd/dns-probe@.service.in
# DNS probe daemon
# Service argument is used as network interface to process traffic from
[Unit]
Description=Daemon for DNS probe
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_BINDIR@/dp-@BACKEND@ -i %i

[Install]
WantedBy=multi-user.target