# systemd/dns-probe.service.in
# DNS probe daemon
# Modify probe's parameters in ExecStart= to fit your use case
[Unit]
Description=Daemon for DNS probe
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_BINDIR@/dns-probe-@BACKEND@

[Install]
WantedBy=multi-user.target
